# dot_config/systemd/user/rog-control-center.service
[Unit]
Description=ROG Control Center (headless)
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/rog-control-center
Restart=always
RestartSec=5
Environment=DISPLAY=:0

[Install]
WantedBy=default.target

